# ansible/roles/logging/templates/journald-backup.service.j2
[Unit]
Description=Upload journald archives via rclone
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/journald-backup.sh
