# ansible/roles/logging/templates/vpn-log-exporter.service.j2
[Unit]
Description=Collect anonymised VPN telemetry
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn-log-exporter.py
