# ansible/roles/logging/templates/vpn-log-exporter.timer.j2
[Unit]
Description=Run VPN telemetry collection every five minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=vpn-log-exporter.service

[Install]
WantedBy=timers.target
